# ansible/files/etc_init_mesos-master.conf
description "mesos master"

start on runlevel [2345]
respawn

exec /usr/bin/mesos-init-wrapper master

