# ansible/files/etc_init_mesos-slave.conf
description "mesos slave"

start on runlevel [2345]
respawn

exec /usr/bin/mesos-init-wrapper slave

